# sparse/lib/systemd/system/pvrinit.service
[Unit]
Requires=new dependency
After=new dependency

[Service]
ExecStart=
ExecStart=new command
